# doc/poison.conf
;supervisor 管理的进程，还需要修改/etc/supervisor/supervisord.conf
;  [supervisord]
;   logfile=/var/log/supervisor/supervisord.log    ; (main log file;default $CWD/supervisord.log)
;   pidfile=/var/run/supervisord.pid               ; (supervisord pidfile;default supervisord.pid)
;   childlogdir=/var/log/supervisor                ; ('AUTO' child log dir, default $TEMP)
;
;    ; 加上下面这两行, 解决 too many open files
;    minfds=65535                   ; min. avail startup file descriptors; default 1024
;    minprocs=65535                 ; min. avail process descriptors;default 200
[program:poison_rpc]
command=poison rpc
stopasgroup=true
killasgroup=true
stopwaitsecs=30
user=root ; User to run as
directory=/home/bolean/
stdout_logfile=/var/log/supervisor/poison_rpc.log ; Where to write log messages
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
autostart=true
autorestart=true
redirect_stderr=true ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,HOME="/mnt/bolean",USER="bolean"
priority=10
[program:poison_server]
command=poison server
stopasgroup=true
killasgroup=true
stopwaitsecs=30
user=root ; User to run as
directory=/home/bolean/
stdout_logfile=/var/log/supervisor/poison_server.log ; Where to write log messages
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
autostart=true
autorestart=true
redirect_stderr=true ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,HOME="/mnt/bolean",USER="bolean"
priority=10
